# scraper/crontab
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

59 23 * * * . /usr/project_env.sh; cd /usr/app && npm run prod >> /var/log/cron.log 2>&1
