# scraper/crontab
* * * * * npm run prod --prefix /usr/app/ >> /var/log/cron.log 2>&1
