# scraper/crontab
* * * * * /usr/local/bin/node /usr/app/src/index.js >> /var/log/cron.log 2>&1
